# supervisor.conf
[program:http]
command=/usr/sbin/httpd -DFOREGROUND
